# pi-server/tv-remote.service
[Unit]
Description=TV Remote Control Server
After=network.target

[Service]
Type=simple
User=kcliang
WorkingDirectory=/home/kcliang/controler/tv-remote-control/pi-server
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/kcliang/.Xauthority
ExecStart=/home/kcliang/controler/tv-remote-control/pi-server/venv/bin/python server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target